# extras/autostart.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/autostart.sh

[Install]
WantedBy=multi-user.target
